# rc.local
# append the following line to the /etc/rc.local
/path/to/startup_script/miner_launcher.sh 15
